# overlays/keystone/board/sifive/hifive-unmatched/unmatched-systemd-units/led-bootstate-green.timer
[Unit]
Description=timer to trun-on green LED
After=getty.target

[Timer]
OnActiveSec=3min
Unit=led-bootstate-green.service

[Install]
WantedBy=timers.target
